# molecule/common/files/fixture-multi.path
[Unit]
Description=Fixture Multi-Label Path

[Path]
PathExists=/tmp/fixture-multi-path-watch
Unit=fixture-multi-a.service

[X-Labels]
Project=myapp
Domain=foo
Env=prod

[Install]
WantedBy=multi-user.target
